# systemd/xrenew.service
[Unit]
Description=Xserver VPS auto extension

[Service]
Type=oneshot
ExecStart={{EXEC_PATH}} extend --auto
