# systemd/xrenew.timer
[Unit]
Description=Xserver VPS auto extension timer

[Timer]
OnCalendar=*-*-* 00,12:00:00
Persistent=true

[Install]
WantedBy=timers.target
